# roles/prometheus_server/templates/PrometheusServer.ini
[program:Prometheus]
command=/usr/local/prometheus/prometheus -config.file=/usr/local/prometheus/prometheus.yml
user=root
autorestart=true  ; auto build when process down
stdout_logfile=/var/log/supervisor/jobs/prometheus_server-supervisord.log ; stdout log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true  ; redirece error output to stderr
