# roles/prometheus_client/templates/PrometheusClient.ini
[program:Prometheus]
command=/usr/local/client_golang/examples/random/random -listen-address=:8080
user=vagrant
autorestart=true  ; auto build when process down
stdout_logfile=/var/log/supervisor/jobs/prometheus_client-supervisord.log ; stdout log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true  ; redirece error output to stderr
